# test/build_sparse_log.cron
#!/bin/bash
# To subscribe to these notices enter your email below and check-in the change.
build=`git rev-parse --short HEAD`
grep "Sparse check failed\!" sparse.*.log
if [ $? -eq 0 ]; then
        mailx -S smtp=smtp://smtp.intel.com -s "[CRON_REG] Sparse check of build $build FAILED" dennis.dalessandro@intel.com < sparse.*.log
        mailx -S smtp=smtp://smtp.intel.com -s "[CRON_REG] Sparse check of build $build FAILED" jubin.john@intel.com < sparse.*.log
else
        mailx -S smtp=smtp://smtp.intel.com -s "[CRON_REG] Sparse check of build $build SUCCESS" dennis.dalessandro@intel.com < sparse.*.log
fi
