# test/daily_checkkconfigsymbols.cron
#!/bin/bash

# Meant to be run by cron but can be run manually.

# Assumes driver source is one dir level up
pushd ../ 1> /dev/null

# Check if checkkymbols.py found
if [ ! -e test/tests/checkkconfigsymbols.py ] ; then
	echo "error: checkkconfigsymbols.py not found!"
	popd 1> /dev/null
	exit 1
fi

# Run test
./test/tests/checkkconfigsymbols.py  -i '(test|scripts|configs)' | tee ./test/tests/checkkconfigsymbols.log

# Process results
popd 1> /dev/null
if ! diff -c tests/checkkconfigsymbols.log tests/checkkconfigsymbols.whitelist; then
	echo "Failed checkkconfigsymbols whitelist comparison"
	exit 1
else
	echo "checkkconfigsymbols test passed"
fi

# If success, remove temporary build log file
rm -f tests/checkkconfigsymbols.log

