# test/daily_update.cron
#!/bin/bash

# Meant to be run by cron but can be run manually.
git remote update
if [ $? -ne 0 ]; then
	echo "Failure: could not update remote"
	exit 1
fi

git pull --rebase
if [ $? -ne 0 ]; then
	echo "Failure: Could not do a pull"
	exit 1
fi

