# test/sperf-05_06_log.cron
#!/bin/bash

# Scrape the build and fpga test logs and mail the results

# To subscribe to these notices add your email below and commit.

build=`git rev-parse --short HEAD`

echo "Nightly regression test for build of $build" > mail.tmp
echo "" > mail.tmp
echo "Log dir: /nfs/sc/disks/fabric_mirror/scratch/wfr_reg_test/daily_logs/RHEL7" >> mail.tmp
echo "Logs: " >> mail.tmp
ls *.log >> mail.tmp
echo "" >> mail.tmp
echo "sperf-05 - 06" >> mail.tmp
echo "-------------" >> mail.tmp
grep -A 100 "Test Summary" sperf-05_06.*.log >> mail.tmp

grep "FAIL" mail.tmp
if [ $? -eq 0 ]; then
	echo "Sending failure mail"
	mailx -S smtp=smtp://smtp.intel.com -s "[CRON_REG] Test on sperf-05_06 for build $build FAILED" dennis.dalessandro@intel.com < mail.tmp
else
	echo "Sending success mail"
	mailx -S smtp=smtp://smtp.intel.com -s "[CRON_REG] Test on sperf-05_06 for build $build SUCCESS" dennis.dalessandro@intel.com < mail.tmp
fi

rm mail.tmp



