# drivers/infiniband/hw/hfi1/test/build_log.cron
#!/bin/bash
# To subscribe to these notices enter your email below and check-in the change.
#Scrape the build log
build=`git rev-parse --short HEAD`
grep "Failed" build.*.log
if [ $? -eq 0 ]; then
        mailx -S smtp=smtp://smtp.intel.com -s "[CRON_REG] Build of $build FAILED" dennis.dalessandro@intel.com < build.*.log
else
        mailx -S smtp=smtp://smtp.intel.com -s "[CRON_REG] Build of $build SUCCESS" dennis.dalessandro@intel.com < build.*.log
fi
