# drivers/infiniband/hw/hfi1/test/daily_reg.cron
#!/bin/bash

#Simple script to be run by cron to do FPGA regression runs

mode=$1
host1=$2
host2=$3

uptime=0

get_uptime () {
	host=$1
	ssh root@$host "cat /proc/uptime" | awk '{print $1}' | awk 'BEGIN { FS="." } {print $1}'
}

power_cycle() {
	host=$1
	echo "Power cycling $host"
	ssh fabserv2 "/nfs/site/proj/fabric/apps/power/bin/power -r $host"
}

check_up() {
	host=$1
	uptime=$2

	attempts=0
	while [ $attempts -lt 10 ]; do
		# Try a ping
		ping -c 1 -t 3 $host
		if [ $? -eq 0 ]; then
			echo "$host pings"
			cur_uptime=$(get_uptime $host)
			if [ $cur_uptime -eq 0 ]; then
				echo "Cur uptime is 0 host may be booting still"
			fi

			if [ $cur_uptime -gt $uptime ]; then
				echo "Host never rebooted after 120 seconds"
				return 1
			fi

			if [ $cur_uptime -lt $uptime ]; then
				echo "Host has rebooted"
				return 0
			fi
		else
			echo "$host does not ping"
		fi
		echo "Checking again in 30 seconds"
		sleep 30
	done 

	return 1
}

do_reboot() {
	host=$1
	uptime=$(get_uptime $host)
	echo "$host has been up for $uptime"
	ssh root@$host 'reboot'&
	echo "Waiting 2 minutes to give host time to shutdown and start booting"
	sleep 120
	check_up $host $uptime
	if [ $? -eq 0 ]; then
		return 0
	else
		echo "$host never came up"
		power_cycle $host1
	
	fi
	check_up $host $uptime
	if [ $? -eq 0]; then
		return 0
	else
		echo "$host never came back after power cycle"
		return 1
	fi
}

# Reboot the hosts to get rid of anything that is fouled up
echo "Rebooting $host1"
do_reboot $host1
if [ $? -ne 0 ]; then
	echo "$host1 never rebooted cleanlyi"
	exit 1
fi

echo "Rebooting $host2"
do_reboot $host2
if [ $? -ne 0 ]; then
	echo "$host2 never rebooted cleanly"
	exit 1
fi
./harness.py --type $mode --nodelist $host1,$host2

# Reboot and hope for the best no need to wait around.
ssh root@$host1 'reboot'&
ssh root@$host2 'reboot'&

